# LJ/listings/137/crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.4797 installed on Fri Apr 15 22:58:22 2005)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
SHELL=/bin/bash
MAILTO=phil

# Catch prairie home companion 11AM sundays
59 10 * * sun /home/phil/catchphc

# Catch Hour of the Wolf 5AM Saturdays 
59 4 * * sat /home/phil/catchthewolf

# Catch Car Talk 10AM Saturdays 
59 9 * * sat /home/phil/catchct

# Catch Schickle Mix 11PM Saturdays
59 21 * * sat /home/phil/catchsch

# Catch Selected Shorts 4PM sundays
59 15 * * sun /home/phil/catchsel

# Catch The Big Broadcast 7PM sundays
59 18 * * sun /home/phil/catchbb




