# LJ/listings/056/crontab
0-55/5 * * * * exec /var/log/sar/bin/sa
0 0 * * * exec /usr/bin/find /var/log/sar/data/ -mtime +14 -exec /bin/rm -f {} \;

